# nemini.service
[Unit]
Description=Nemini - The simple Gemini server service

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/usr/local/bin/nemini
User=nemini

[Install]
WantedBy=multi-user.target
